# luci-app-timedtask/root/etc/uci-defaults/luci-timedtask
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@timedtask[-1]
	add ucitrack timedtask
	set ucitrack.@timedtask[-1].init=timedtask
	commit ucitrack
EOF

rm -f /tmp/luci-*
exit 0
