# luci-app-pwdHackDeny/root/etc/init.d/pwdHackDeny
#!/bin/sh /etc/rc.common
# wulishui 20210314 v6.0.4
# Author: wulishui <wulishui@gmail.com>

START=92
USE_PROCD=1
tmp_dir="/tmp/pwdHackDeny"
deny_dir="/etc/pwdHackDeny"

start_instance() {
	procd_open_instance
	procd_set_param command /etc/pwdHackDeny.sh
	procd_set_param respawn
	procd_set_param stderr 1
	procd_close_instance
}

chk_logfile() {
	echo "## 每行一条，可手动添加或删除，内网MAC或IPV4/IPV4段(如192.168.18.10-20)、外网IPV4/IPV6，或包含此内容的文本、日志。内网MAC外网IP行首加#将不拦截。Author: wulishui" >$tmp_dir/$1.log.bak
	grep -v '^## ' /etc/$1.log | grep '^[#/!]' >>$tmp_dir/$1.log.bak && echo "----------- 以上客户端不拦截 -----------" >>$tmp_dir/$1.log.bak
	sshlog=$(cat /etc/$1.log 2>/dev/null | sed '/^----/d' | sed '/^'[#!/]'/d' | sed 's/  */ /g' | sed 's/^ *//g' | sort -n | uniq -i | awk NF)
	if [ -n "$sshlog" ]; then
		echo "$sshlog" | while read line; do
			MAC=$(echo "$line" | egrep -o "([A-Fa-f0-9]{2}[:-]){5}[A-Fa-f0-9]{2}" | sed 's/-/:/g' | head -1)
			[ -z "$MAC" ] && IPR=$(echo "$line" | egrep -o "(([0-9]|[1-9][0-9]|1[0-9]{2}|2[0-4][0-9]|25[0-5])\.){3}([0-9]|[1-9][0-9]|1[0-9]{2}|2[0-4][0-9]|25[0-5])-([1-9]|[1-9][0-9]|1[0-9][0-9]|2[0-5][0-4])$" | head -1) && IPRANGE=$(echo "$IPR" | awk -F '[.-]' '{if ($5=="") print $1"."$2"."$3"."$4 ; else if ($8=="" && $4<$5) print $1"."$2"."$3"."$4"-"$1"."$2"."$3"."$5 ; else if ($8=="" && $4>$5) print $1"."$2"."$3"."$5"-"$1"."$2"."$3"."$4 ; else if ($8>$4) print $1"."$2"."$3"."$4"-"$1"."$2"."$3"."$8 ; else if ($8<$4) print $1"."$2"."$3"."$8"-"$1"."$2"."$3"."$4  }' | awk -F '[.-]' '{if ($8=="") print $1"."$2"."$3"."$4 ; else if ($8<255) print $1"."$2"."$3"."$4"-"$1"."$2"."$3"."$8 ; else if ($8>254) print $1"."$2"."$3"."$4"-"$1"."$2"."$3".""254"}')
			[ -z "$MAC" -a -z "$IPR" ] && IP4=$(echo "$line" | egrep -o "(([0-9]|[1-9][0-9]|1[0-9]{2}|2[0-4][0-9]|25[0-5])\.){3}([0-9]|[1-9][0-9]|1[0-9]{2}|2[0-4][0-9]|25[0-5])" | head -1)
			[ -z "$MAC" -a -z "$IPR" -a -z "$IP4" ] && IP6=$(echo "$line" | egrep -o "(s*((([0-9A-Fa-f]{1,4}:){7}([0-9A-Fa-f]{1,4}|:))|(([0-9A-Fa-f]{1,4}:){6}(:[0-9A-Fa-f]{1,4}|((25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)(\.(25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)){3})|:))|(([0-9A-Fa-f]{1,4}:){5}(((:[0-9A-Fa-f]{1,4}){1,2})|:((25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)(\.(25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)){3})|:))|(([0-9A-Fa-f]{1,4}:){4}(((:[0-9A-Fa-f]{1,4}){1,3})|((:[0-9A-Fa-f]{1,4})?:((25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)(\.(25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)){3}))|:))|(([0-9A-Fa-f]{1,4}:){3}(((:[0-9A-Fa-f]{1,4}){1,4})|((:[0-9A-Fa-f]{1,4}){0,2}:((25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)(\.(25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)){3}))|:))|(([0-9A-Fa-f]{1,4}:){2}(((:[0-9A-Fa-f]{1,4}){1,5})|((:[0-9A-Fa-f]{1,4}){0,3}:((25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)(\.(25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)){3}))|:))|(([0-9A-Fa-f]{1,4}:){1}(((:[0-9A-Fa-f]{1,4}){1,6})|((:[0-9A-Fa-f]{1,4}){0,4}:((25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)(\.(25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)){3}))|:))|(:(((:[0-9A-Fa-f]{1,4}){1,7})|((:[0-9A-Fa-f]{1,4}){0,5}:((25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)(\.(25[0-5]|2[0-4]\d|1\d\d|[1-9]?\d)){3}))|:)))(%.+)?\s*)")
			if [ -n "$MAC" ]; then
				[ -z "$(grep -w "$MAC" $tmp_dir/$1.log.bak 2>/dev/null)" ] && {
					echo "$MAC" >>$tmp_dir/$1.log.bak
					iptables -w -A $2 -m mac --mac-source $MAC -j DROP 2>/dev/null
					ip6tables -w -A $2 -m mac --mac-source $MAC -j DROP 2>/dev/null
				}
			elif [ -n "$IPR" ]; then
				[ -z "$(grep -w "$IPR" $tmp_dir/$1.log.bak 2>/dev/null)" ] && {
					echo "$IPR" >>$tmp_dir/$1.log.bak
					iptables -w -A $2 -m iprange --src-range $IPRANGE -j DROP 2>/dev/null
				}
			elif [ -n "$IP4" ]; then
				[ -z "$(grep -w "$IP4" $tmp_dir/$1.log.bak 2>/dev/null)" ] && {
					echo "$IP4" >>$tmp_dir/$1.log.bak
					iptables -w -A $2 -s $IP4 -j DROP 2>/dev/null
				}
			elif [ -n "$IP6" ]; then
				[ -z "$(grep -w "$IP6" $tmp_dir/$1.log.bak 2>/dev/null)" ] && {
					echo "$IP6" >>$tmp_dir/$1.log.bak
					ip6tables -w -A $2 -s $IP6 -j DROP 2>/dev/null
				}
			fi
			unset IPRANGE
		done
	fi
	[ -s $tmp_dir/$1.log.bak ] && cat $tmp_dir/$1.log.bak >/etc/$1.log
}

start_service() {
	kill -9 $(pgrep 'pwdHackDeny.sh') >/dev/null 2>&1
	sh /tmp/log/pwdHackDeny.ipt.cls 2>/dev/null
	enabled=$(uci -q get pwdHackDeny.pwdHackDeny.enabled)
	[ "$enabled" != 1 -a -d $deny_dir ] && touch $deny_dir/logclear
	[ "$enabled" != 1 ] && return 1
	if [ -e $deny_dir/logclear ]; then
		[ -s $deny_dir/badip.log.web ] && grep -v '^----' $deny_dir/badip.log.web >>$deny_dir/bak.log.web && \
		echo "--------"$(date +"%Y-%m-%d %H:%M:%S")" : 日志已移动到 $deny_dir/bak.log.web .--------" >$deny_dir/badip.log.web
		[ -s $deny_dir/badip.log.ssh ] && grep -v '^----' $deny_dir/badip.log.ssh >>$deny_dir/bak.log.ssh && \
		echo "--------"$(date +"%Y-%m-%d %H:%M:%S")" : 日志已移动到 $deny_dir/bak.log.ssh .--------" >$deny_dir/badip.log.ssh
		rm -f $deny_dir/logclear 2>/dev/null
	fi
	mkdir $deny_dir $tmp_dir 2>/dev/null
	touch $deny_dir/badip.log.web $deny_dir/badip.log.ssh $deny_dir/badhosts.web $deny_dir/badhosts.ssh

	redirect_https=$(uci -q get uhttpd.main.redirect_https)
	port1=$(uci -q get dropbear.@dropbear[0].Port) || port1=22
	port2=$(uci -q get uhttpd.main.listen_http | egrep -o "(:[1-9][0-9]{0,3}|[1-5][0-9]{4}|6[0-4][0-9]{3}|65[0-4][0-9]{2}|655[0-2][0-9]{1}|6553[0-5])" | sed 's/://g' | sort -n | uniq -i | awk NF) || port2=80
	port2s=$(uci -q get uhttpd.main.listen_https | egrep -o "(:[1-9][0-9]{0,3}|[1-5][0-9]{4}|6[0-4][0-9]{3}|65[0-4][0-9]{2}|655[0-2][0-9]{1}|6553[0-5])" | sed 's/://g' | sort -n | uniq -i | awk NF) || port2s=443
	#configip4=$(uci -q get network.lan.ipaddr)
	#configip6=$(uci -q get network.globals.ula_prefix)

	iptables -w -N pwdHackDenySSH 2>/dev/null || iptables -w -F pwdHackDenySSH 2>/dev/null
	ip6tables -w -N pwdHackDenySSH 2>/dev/null || ip6tables -w -F pwdHackDenySSH 2>/dev/null
	iptables -w -N pwdHackDenyWEB 2>/dev/null || iptables -w -F pwdHackDenyWEB 2>/dev/null
	ip6tables -w -N pwdHackDenyWEB 2>/dev/null || ip6tables -w -F pwdHackDenyWEB 2>/dev/null

	iptables -w -C INPUT -p tcp --dport "$port1" -j pwdHackDenySSH 2>/dev/null || iptables -w -I INPUT -p tcp --dport "$port1" -j pwdHackDenySSH 2>/dev/null
	ip6tables -w -C INPUT -p tcp --dport "$port1" -j pwdHackDenySSH 2>/dev/null || ip6tables -w -I INPUT -p tcp --dport "$port1" -j pwdHackDenySSH 2>/dev/null

	echo "$port2" | while read line; do
		iptables -w -C INPUT -p tcp --dport "$line" -j pwdHackDenyWEB 2>/dev/null || iptables -w -I INPUT -p tcp --dport "$line" -j pwdHackDenyWEB 2>/dev/null
		ip6tables -w -C INPUT -p tcp --dport "$line" -j pwdHackDenyWEB 2>/dev/null || ip6tables -w -I INPUT -p tcp --dport "$line" -j pwdHackDenyWEB 2>/dev/null
	done

	if [ "$redirect_https" == 1 ]; then
		echo "$port2s" | while read line; do
			iptables -w -C INPUT -p tcp --dport "$line" -j pwdHackDenyWEB 2>/dev/null || iptables -w -I INPUT -p tcp --dport "$line" -j pwdHackDenyWEB 2>/dev/null
			ip6tables -w -C INPUT -p tcp --dport "$line" -j pwdHackDenyWEB 2>/dev/null || ip6tables -w -I INPUT -p tcp --dport "$line" -j pwdHackDenyWEB 2>/dev/null
		done
	fi

	chk_logfile SSHbadip pwdHackDenySSH
	chk_logfile WEBbadip pwdHackDenyWEB

	cat >/tmp/log/pwdHackDeny.ipt.cls <<-EOF
	iptables  -w -D INPUT -p tcp --dport "$port1"  -j pwdHackDenySSH 2>/dev/null
	ip6tables -w -D INPUT -p tcp --dport "$port1"  -j pwdHackDenySSH 2>/dev/null
	iptables  -w -D INPUT -p tcp --dport "$PTV4"  -j pwdHackDenyWEB 2>/dev/null
	ip6tables -w -D INPUT -p tcp --dport "$PTV6"  -j pwdHackDenyWEB 2>/dev/null
	iptables  -w -D INPUT -p tcp --dport "$PTSV4" -j pwdHackDenyWEB 2>/dev/null
	ip6tables -w -D INPUT -p tcp --dport "$PTSV6" -j pwdHackDenyWEB 2>/dev/null
	iptables  -w -F pwdHackDenySSH 2>/dev/null
	iptables  -w -X pwdHackDenySSH 2>/dev/null
	ip6tables -w -F pwdHackDenySSH 2>/dev/null
	ip6tables -w -X pwdHackDenySSH 2>/dev/null
	iptables  -w -F pwdHackDenyWEB 2>/dev/null
	iptables  -w -X pwdHackDenyWEB 2>/dev/null
	ip6tables -w -F pwdHackDenyWEB 2>/dev/null
	ip6tables -w -X pwdHackDenyWEB 2>/dev/null
	rm -f /tmp/log/pwdHackDeny.ipt.cls 2>/dev/null
	EOF
	rm -f $tmp_dir/*.log.bak
	touch $tmp_dir/syslog
	start_instance
}

service_triggers() {
	procd_add_reload_trigger 'pwdHackDeny'
}

stop_service() {
	sh /tmp/log/pwdHackDeny.ipt.cls 2>/dev/null
	# rm -rf $tmp_dir 2>/dev/null
}

restart() {
	stop_service
	start_service
}
