# luci-app-deluge/root/etc/uci-defaults/luci-app-deluge
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@deluge[-1]
	add ucitrack deluge
	set ucitrack.@deluge[-1].init=deluge
	commit ucitrack
EOF

rm -rf /tmp/luci-*
exit 0
