# luci-app-wizard/root/etc/init.d/wizard
#!/bin/sh /etc/rc.common

START=19

extra_command "reconfig" ""

boot() {
	XBOOT=1 start
}

add_wizard() {
	[ "x$XBOOT" = "x1" ] && return 0
	get_config="dhcp ipv6 lan_dns lan_gateway lan_ipaddr lan_netmask old_ipv6 old_siderouter old_wifi_key old_wifi_ssid siderouter wan_pppoe_pass wan_pppoe_user wan_proto wifi_key wifi_ssid"
	for x in $get_config; do config_get $x $1 $x; done

	case $wan_proto in
		dhcp)
			[ "$(uci -q get network.wan.proto)" != dhcp ] && {
				uci -q set network.wan.proto='dhcp'
				uci -q del network.wan.username
				uci -q del network.wan.password
			}
		;;
		pppoe)
			[ "$wan_pppoe_user" != "$(uci -q get network.wan.username)" -o "$wan_pppoe_pass" != "$(uci -q get network.wan.password)" ] && {
				uci -q set network.wan.proto='pppoe'
				uci -q set network.wan.username="$wan_pppoe_user"
				uci -q set network.wan.password="$wan_pppoe_pass"
			}
		;;
	esac

	if [ "$siderouter" != "$old_siderouter" ]; then
		if [ "$lan_gateway" != "$(uci -q get network.lan.gateway)" -o  "$dhcp" != "$(uci -q get dhcp.lan.ignore)" ]; then
			if [ "$siderouter" = 1 ]; then
				uci -q set network.lan.gateway="$lan_gateway"
				[ "$dhcp" = 1 ] && uci -q set dhcp.lan.ignore="1" || uci -q del dhcp.lan.ignore
				uci -q set firewall.@zone[0].masq='1'
			else
				uci -q del dhcp.lan.ignore
				uci -q del network.lan.dns
				uci -q del network.lan.gateway
				uci -q del firewall.@zone[0].masq
			fi
		fi
		uci -q set wizard.default.old_siderouter="$siderouter"
	fi
	[ "$lan_dns" != "$(uci -q get network.lan.dns)" ] && uci -q set network.lan.dns="$lan_dns"
	[ "$lan_ipaddr" != "$(uci -q get network.lan.ipaddr)" ] && uci -q set network.lan.ipaddr="$lan_ipaddr"
	[ "$lan_netmask" != "$(uci -q get network.lan.netmask)" ] && uci -q set network.lan.netmask="$lan_netmask"

	if [ "$ipv6" != "$old_ipv6" ]; then
		case "$ipv6" in
		0)
			uci -q set network.wan.ipv6='0'
			uci -q set network.wan.delegate='0'
			;;
		1)
			uci -q delete dhcp.lan.ra
			uci -q delete dhcp.lan.dhcpv6
			uci -q set network.wan.ipv6='0'
			uci -q set network.lan.delegate='0'
			uci -q set network.wan.delegate='0'
			uci -q delete network.globals.ula_prefix
			;;
		2)
			uci -q del network.lan.delegate
			uci -q del network.wan.delegate
			uci -q set dhcp.lan.ra='hybrid'
			uci -q set network.wan.ipv6='auto'
			uci -q set dhcp.lan.dhcpv6='hybrid'
			;;
		esac
		uci -q set wizard.default.old_ipv6="$ipv6"
	fi

	test -s /etc/config/wireless && {
		wifi_setup_radio() {
			local radio=$1
			uci -q get wireless.${radio} && {
				if [ "$(uci -q get wireless.${radio}.band)" = "5g" ]; then
					uci -q set wireless.default_${radio}.ssid="${wifi_ssid}_5G"
				else
					uci -q set wireless.default_${radio}.ssid="${wifi_ssid}_2.4G"
				fi
				uci -q set wireless.default_${radio}.device="$radio"
				if [ "$wifi_key" ]; then
					uci -q set wireless.default_${radio}.key="$wifi_key"
					uci -q set wireless.default_${radio}.encryption='psk2'
				else
					uci -q set wireless.default_${radio}.encryption='none'
				fi
			}
		}

		if [ "$wifi_ssid" != "$old_wifi_ssid" -o "$wifi_key" != "$old_wifi_key" ]; then
			test -n "$wifi_ssid" && {
				for radio in radio0 radio1 radio2 radio3; do
					wifi_setup_radio "$radio"
				done
				uci commit wireless
			}
			uci -q set wizard.default.old_wifi_key="$wifi_key"
			uci -q set wizard.default.old_wifi_ssid="$wifi_ssid"
		fi
	}

	uci commit dhcp
	uci commit wizard
	uci commit network
	uci commit firewall
	/etc/init.d/network reload &
	/etc/init.d/dnsmasq reload &
}

reconfig() {
	uci -q set wizard.default.lan_dns="$(uci -q get network.lan.dns)"
	uci -q set wizard.default.wan_proto="$(uci -q get network.wan.proto)"
	uci -q set wizard.default.lan_ipaddr="$(uci -q get network.lan.ipaddr)"
	uci -q set wizard.default.lan_netmask="$(uci -q get network.lan.netmask)"
	uci -q set wizard.default.lan_gateway="$(uci -q get network.lan.gateway)"
	uci -q set wizard.default.wan_pppoe_user="$(uci -q get network.wan.username)"
	uci -q set wizard.default.wan_pppoe_pass="$(uci -q get network.wan.password)"
	[ -n "$(uci -q get dhcp.lan.ignore)" ] && uci -q set wizard.default.dhcp="0" || uci -q del wizard.default.dhcp
	uci commit wizard
}

start() {
	config_load wizard
	config_foreach add_wizard wizard
}

restart() {
	XRELOAD=1 start
}
