# luci-app-deluge/root/etc/init.d/deluge
#!/bin/sh /etc/rc.common

START=99

start() {
	config_load deluge
	get_config="user profile_dir download_dir log_dir Locale password port enabled"

	for rt in $get_config; do
		config_get_bool $rt main $rt
		config_get $rt main $rt
	done
	[ $port ] && port="-p $port"
	[ $user ] && user="-u $profile_dir"
	[ $profile_dir ] && profile_dir="-c $profile_dir"
	[ $log_dir ] && log_dir="-l $log_dir/deluge.log"
	[ $enabled ] && {
		/usr/bin/python3 /usr/bin/deluged $profile_dir $user $log_dir -L error
		/usr/bin/python3 /usr/bin/deluge-web $profile_dir $port $log_dir -L error
	}
}

stop() {
	# killall deluged
	# killall deluge-web
	kill $(ps -w | grep -w /usr/bin/deluged | grep -v grep | awk '{print $1}') >/dev/null 2>&1
	kill $(ps -w | grep -w deluge-web | grep -v grep | awk '{print $1}') >/dev/null 2>&1
}
