# luci-app-ddnsto/root/etc/init.d/ddnsto
#!/bin/sh /etc/rc.common

START=99
token=$(uci get ddnsto.@global[0].token)
enable=$(uci get ddnsto.@global[0].enable)

start() {
	[[ $enable == 0 ]] && exit 0
	service_start /usr/bin/ddnsto -u $token -d >/dev/null 2>&1 &
	logger -t ddnsto "Start ddnsto service"
}

stop() {
	killall ddnsto 2>/dev/null
	logger -t ddnsto "Stop ddnsto service"
}

boot() {
	delay_start=$(uci get ddnsto.@global[0].start_delay || echo 0)
	[[ $delay_start -gt 0 && $enable -eq 1 ]] && {
		logger -t ddnsto "延时$delay_start秒启动 ddnsto"
		sleep "$delay_start"
	}
	start
}
