# luci-app-ddnsto/root/etc/init.d/ddnsto
#!/bin/sh /etc/rc.common

START=99
token=$(uci_get ddnsto ddnsto token)
enable=$(uci_get ddnsto ddnsto enable)

start() {
    [ x$enable = x1 ] || return 0
    [ -z $token ] && {
        logger -t ddnsto "token值为空"
        return 1
    }
    /usr/bin/ddnsto -u $token -d >/dev/null 2>&1 && \
    logger -t ddnsto "启动 ddnsto 服务"
}

stop() {
    pidof ddnsto >/dev/null && {
        killall ddnsto 2>/dev/null && \
        logger -t ddnsto "关闭 ddnsto 服务"
    }
}

boot() {
    [ x$enable = x1 ] || return 0
    logger -t ddnsto "延时 $delay_start 秒启动 ddnsto"
    sleep $(uci_get ddnsto ddnsto start_delay)
    start
}
